# cron-jobs/all.cron
54 * * * * /var/www/cron-jobs/update_princetonfood
